# roles/configure_slurm/templates/systemd/munge.service.d/trix_specific.conf.j2
{{ ansible_managed | comment }}

[Unit]
After=remote-fs.target
Requires=remote-fs.target

[Service]
ExecStart=
ExecStart=/usr/sbin/munged  --key-file {{ trix_path_shared }}/etc/munge/munge.key

